# src/crontab
59 1 * * 1 root /usr/sbin/smartschedule startslow-and-report
