# src/crontab
1 2 * * 1 root /usr/sbin/smartschedule start
1 8 * * 1 root /usr/sbin/smartschedule report
